# services/deluged/run
#!/bin/bash

exec /sbin/setuser abc /usr/bin/deluged -c /config -d --loglevel=info -l /config/deluged.log
