# root/etc/services.d/deluge-web/run
#!/usr/bin/with-contenv bash

DELUGE_LOGLEVEL=${DELUGE_LOGLEVEL:-warning}


exec \
    s6-setuidgid abc /usr/bin/deluge-web \
    -d -c /config --loglevel="${DELUGE_LOGLEVEL}"
